# deploy/hermitstash.service
# HermitStash — systemd service for bare metal deployment
#
# Install:
#   sudo cp hermitstash.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now hermitstash
#
# Prerequisites:
#   - Node.js 24+ installed (for PQC / OpenSSL 3.5 support)
#   - Application at /opt/hermitstash (git clone or copy)
#   - User 'hermit' created: sudo useradd -r -m -d /opt/hermitstash -s /bin/false hermit
#   - Data dirs: sudo mkdir -p /opt/hermitstash/{data,uploads}
#   - Permissions: sudo chown -R hermit:hermit /opt/hermitstash
#   - tmpfs mount for in-memory DB (add to /etc/fstab):
#     tmpfs /opt/hermitstash/shm tmpfs size=256M,mode=700,uid=hermit,gid=hermit 0 0

[Unit]
Description=HermitStash — Post-quantum encrypted file uploads
Documentation=https://github.com/dotCooCoo/hermitstash
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hermit
Group=hermit
WorkingDirectory=/opt/hermitstash

# Environment
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HERMITSTASH_TMPDIR=/opt/hermitstash/shm

# Vault passphrase protection (opt-in, v1.9+)
# Default: disabled (vault.key is plaintext JSON, filesystem-protected only).
# To enable:
#   1. Store the passphrase in a secured file owned by root, mode 0400:
#        sudo install -m 0400 -o root -g root /dev/stdin /etc/hermitstash/vault-passphrase
#      then type (or pipe) the passphrase. systemd will read the file at start.
#   2. Stop the service, run `sudo -u hermit node scripts/vault-passphrase-setup.js`
#      to migrate plaintext vault.key → vault.key.sealed.
#   3. Uncomment the two lines below and `systemctl daemon-reload`.
#   4. Start the service.
# LOSS of the passphrase file = LOSS of all encrypted data. Back it up safely.
#
# Environment=VAULT_PASSPHRASE_MODE=required
# Environment=VAULT_PASSPHRASE_FILE=/etc/hermitstash/vault-passphrase

# ── PEM at-rest sealing for CA + TLS keys (v1.9.4+) ──────────────────────
# Tristate: auto (default) | required | disabled. Flip to "required" after
# sealing via admin UI (v1.9.6+) or the CLI scripts.
Environment="CA_KEY_SEALED=auto"
Environment="TLS_KEY_SEALED=auto"

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
ReadWritePaths=/opt/hermitstash/data /opt/hermitstash/uploads /opt/hermitstash/shm

# Process
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hermitstash

[Install]
WantedBy=multi-user.target
